# configuration/geod.service
[Unit]
Description=Seismic signal distribution daemon

[Service]
ExecStart=/usr/bin/geod

[Install]
WantedBy=multi-user.target